# docker-web-v2/crontab
10 */12 * * * docker restart end-gfw-web
